# jsx_machine/flask_service.service
[Unit]
Description=Job that runs the python SimpleHTTPServer daemon
Documentation=man:SimpleHTTPServer(1)

[Service]
Type=simple
WorkingDirectory=/home
ExecStart=/usr/bin/python3.6 /mnt/data/scripts/flask_server.py &
ExecStop=/bin/kill `/bin/ps aux | /bin/grep /mnt/data/scripts/flask_server.py | /bin/grep -v grep | /usr/bin/awk '{ print $2 }'`
Restart=on-abort

[Install]
WantedBy=multi-user.target